# resetbutton.service
[Unit]
Description=ResetButtonService

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /home/pi/shared/rotary-box-rpi/resetbutton.py
After=gitpull.service

[Install]
WantedBy=multi-user.target
