# rotarymeds.service
[Unit]
Description=RotaryMedsService

[Service]
Type=simple
User=pi
ExecStart=/home/pi/rotarymeds.sh

[Install]
WantedBy=multi-user.target
